# launchd/com.higgins.app.plist.tpl
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.higgins.app</string>
  <key>ProgramArguments</key>
  <array>
    <string>{{NODE_PATH}}</string>
    <string>{{HIGGINS_ROOT}}/index.mjs</string>
  </array>
  <key>WorkingDirectory</key>
  <string>{{HIGGINS_ROOT}}</string>
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>
  <key>StandardOutPath</key>
  <string>{{HIGGINS_ROOT}}/logs/higgins.log</string>
  <key>StandardErrorPath</key>
  <string>{{HIGGINS_ROOT}}/logs/higgins.err.log</string>
  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
  </dict>
  <key>ProcessType</key>
  <string>Background</string>
</dict>
</plist>
