# launchd/com.higgins.calsync.plist.tpl
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.higgins.calsync</string>
  <key>ProgramArguments</key>
  <array>
    <string>{{NODE_PATH}}</string>
    <string>{{HIGGINS_ROOT}}/skills/calendar/cal-sync.mjs</string>
  </array>
  <key>WorkingDirectory</key>
  <string>{{HIGGINS_ROOT}}</string>
  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key>
    <integer>6</integer>
    <key>Minute</key>
    <integer>0</integer>
  </dict>
  <key>StandardOutPath</key>
  <string>{{HIGGINS_ROOT}}/logs/calsync.log</string>
  <key>StandardErrorPath</key>
  <string>{{HIGGINS_ROOT}}/logs/calsync.err.log</string>
  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
  </dict>
</dict>
</plist>
